# puppet/files/udev/mount.rules
SUBSYSTEM=="bdi",ACTION=="add",RUN+="/etc/init.d/server start"
SUBSYSTEM=="bdi",ACTION=="remove",RUN+="/etc/init.d/server stop"